# game_server/bin/game_server.service
[Unit]
Description=Game Server

[Service]
ExecStart=/usr/bin/forever --minUptime 1000 --spinSleepTime 1000 /home/sessi0n/game_server/main 00
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=game_server
User=sessi0n
Group=sessi0n
Environment="NODE_ENV=production"
Environment="SERVICE_ENV=rusty"

[Install]
WantedBy=multi-user.target
